# tasks/heartbeat/plist/com.lifeos.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lifeos.heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VAULT_ROOT__/tasks/heartbeat/runner.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__VAULT_ROOT__</string>

    <key>StartInterval</key>
    <integer>__INTERVAL_SECONDS__</integer>

    <key>StandardOutPath</key>
    <string>__VAULT_ROOT__/state/heartbeat/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__VAULT_ROOT__/state/heartbeat/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
